# config/umamao.cron
0 * * * * ~/refresh_related_topics_and_suggestions
